# files/etc/rc.local
#!/bin/sh
logger -t nswa-preact "run auto act"
/lib/ranga/libexec/act --auto
exit 0
